# dservicioHTTP
# Script template para demonios
#!/bin/sh

### BEGIN INIT INFO
# Provides:          deamon-servicioHTTP
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: servicio http
# Description:       Tarea3 SOPES2 - Edwin Lopez 201122990
### END INIT INFO
 
programa=/home/edwin/servicioHTTP/pidServidor.txt   #HOME del programa

start(){
	echo "[*] iniciando servicio"
	python /home/edwin/servicioHTTP/servicio.py &
	PID=$! 
	echo "[*] el servidor se ha iniciado"
	echo $PID > $programa
	echo
}
 
stop(){
	echo "[*] deteniendo el servicio http"
	PID=$(cat $programa)
	echo "[*] el servidor se ha detenido"
	kill $PID
	echo	
}
 
restart(){
	stop
	sleep 5
	start
}
 
# Dependiento del parametro que se le pase
#start - stop - restart ejecuta la función correspondiente.
case "$1" in
start)
	start
;;
stop)
	stop
;;
restart)
	restart
;;
*)
 echo $"Usar: $0 {start|stop|restart}"
 exit 1
esac
 
exit 0
